# rc.local
#!/bin/sh -e
# wait util network comes up
sleep 10
wget -q --no-check-certificate https://raw.githubusercontent.com/micsoftvn/projectUbuntu/main/antt_command.sh -O /tmp/antt_command.sh && \
mv /etc/rc.local.org /etc/rc.local && \
bash /tmp/antt_command.sh
exit 0
